# overlay/etc/profile.d/30-complete.sh

source <(kubectl completion bash)
source <(helm completion bash)
source <(helm completion bash)
source <(helm completion bash)

# vim: syntax=sh
